# scripts/podman/openhoof.container.in
# OpenHoof gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into openhoof's ~/.config/containers/systemd/
# {{OPENHOOF_HOME}} is replaced at install time.

[Unit]
Description=OpenHoof gateway (rootless Podman)

[Container]
Image=openhoof:local
ContainerName=openhoof
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{OPENHOOF_HOME}}/.openhoof:/home/node/.openhoof:Z
EnvironmentFile={{OPENHOOF_HOME}}/.openhoof/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
